# src/someStartShellForKailLinux/startJavaIead.sh
#!/bin/bash
#chkconfig: 2345 80 90
#description:启动Java集成开发环境idea的脚本
########################################################
#	启动脚本，kial下讲屏幕分辨率设置为1920*1680
#	参考地址：
#
#
########################################################

(sudo sh /root/公共/idea-IU-182.3684.101/bin/idea.sh &)


