# src/someStartShellForKailLinux/startJavaIead.sh
#!/bin/bash
#chkconfig: 2345 80 90
#description:启动Java集成开发环境idea的脚本
########################################################
#	启动脚本，kial下讲屏幕分辨率设置为1920*1680
#	参考地址：	http://baijiahao.baidu.com/s?id=1594632413437230067&wfr=spider&for=pc
#				https://blog.csdn.net/daidaineteasy/article/details/54864776
#				https://blog.csdn.net/linuxshine/article/details/50717272
########################################################

(sudo sh /root/公共/idea-IU-182.3684.101/bin/idea.sh &)


